# lexia-backend/s6/s6-rc.d/minio-bucket/up
/command/with-contenv /usr/local/bin/minio-bucket-init.sh
