# lexia-backend/s6/s6-rc.d/auth/run
#!/command/with-contenv sh
exec 2>&1
cd /app
echo "[auth] starting node server on :3000"
exec node /app/dist/server.js
