# lexia-backend/s6/s6-rc.d/auth-migrate/up
/command/with-contenv /usr/local/bin/auth-migrate.sh
